# 2022-06-06-puppet_printers_reset/puppet_printers_reset/files/profile.d/ies-printers
#!/bin/bash
sudo /usr/local/sbin/ies-printers
